# applications/luci-app-i-love-luci/root/etc/uci-defaults/90_luci-app-i-love-luci
#!/bin/sh

rm -rf /tmp/luci-indexcache /tmp/luci-modulecache
/etc/init.d/rpcd reload >/dev/null 2>&1 || true

exit 0
